# server/usr/lib/systemd/system/backup-volumes.timer
[Unit]
Description=Periodically backup docker volumes
[Timer]
OnCalendar=*-*-* 20:00:00
[Install]
WantedBy=timers.target
